# Koan/Config/crontab
#Mins  Hours  Days   Months  Day of the week    Command
0,15,30,45 * * * * /Scripts/newroute.sh # every 15 minutes
0,15,30,45 * * * * /Scripts/randomise_mac.sh # every 15 minutes
